# private_dot_config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Exec=bash ~/.config/autostart/autorun.sh
Hidden=true
