# private_dot_config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Exec=polybar
Hidden=true
